# extra_files/temp.timer
[Unit]
Description=Run tempcheck every minute

[Timer]
OnCalendar=*:0/1
Persistent=true

[Install]
WantedBy=timers.target

